# rpm/istatd-agent.postinst
#!/bin/sh
set -e

if [ "$1" != "1" ]; then
        exit 0
fi

if [ -x "/etc/init.d/istatd-agent" ]; then
        /usr/sbin/chkconfig istatd-agent on
        if [ -x /usr/sbin/service ]; then
                /usr/sbin/service istatd-agent start ||true
        else
                /etc/init.d/istatd-agent start ||true
        fi
fi

exit 0
