# rpm/istatd-agent.prerm
#!/bin/sh
set -e

if [ "$1" != "0" ]; then
        exit 0
fi

if [ -x "/etc/init.d/istatd-agent" ]; then
        /usr/sbin/chkconfig istatd-agent off
        if [ -x /usr/sbin/service ]; then
                /usr/sbin/service istatd-agent stop || true
        else
                /etc/init.d/istatd-agent stop || true
        fi
        /usr/sbin/chkconfig --del istatd-agent
fi

exit 0
