# rpm/istatd-server.prerm
#!/bin/sh
set -e

if [ "$1" != "0" ]; then
        exit 0
fi

if [ -x "/etc/init.d/istatd-server" ]; then
        /usr/sbin/chkconfig istatd-server on || true
        if [ -x /usr/sbin/service ]; then
                /usr/sbin/service istatd-server stop || true
        else
                /usr/sbin/service stop || true
        fi
        /usr/sbin/chkconfig --del istatd-server
fi

exit 0
